# templates/default/sysvinit.erb
#!/bin/sh
#
# /etc/init.d/<%= @name %> -- Startup/shutdown script for Tomcat
#
#
### BEGIN INIT INFO
# Provides: tomcat
# Required-Start: $Local_fs $remote_fs $network $sysLog
# Required-Stop: $Local_fs $remote_fs $network $sysLog
# Should-Start: $named
# Should-Stop $named
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Start and stop the Tomcat servlet application
# Short-Description: Start and stop Tomcat
### END INIT INFO

RUNAS=<%= @user %>
KILL_WAIT=<%= @environment['KILL_WAIT'] || 45 %>
export CATALINA_HOME=<%= @environment['CATALINA_HOME'] || @directory %>
export CATALINA_BASE=<%= @environment['CATALINA_BASE'] || @directory %>
export CATALINA_PID=<%= @pid_file %>
<%- @environment.each do |key, val| -%>
<%- next if %w(CATALINA_HOME CATALINA_BASE CATALINA_PID KILL_WAIT).include? key %>
export <%= key %>="<%= val %>"
<%- end -%>

if [ -z "$CATALINA_HOME" ]; then
    echo "CATALINA_HOME must be set!"
    exit 5
fi

require_root() {
    if [ $(id -u) -ne 0 ]; then
        echo "Root permission is required for this operation"
        exit 4
    fi
}

tomcat_start() {
    require_root
    PID_DIR=$(dirname $CATALINA_PID)
    [ -d "$PID_DIR" ] || mkdir -p $PID_DIR
    touch $CATALINA_PID
    chown $RUNAS: $CATALINA_PID
    /bin/su -p -s /bin/sh -c "$CATALINA_HOME/bin/catalina.sh start" $RUNAS
}

tomcat_stop() {
    require_root
    /bin/su -p -s /bin/sh -c "$CATALINA_HOME/bin/catalina.sh stop $KILL_WAIT -force" $RUNAS
}

tomcat_status() {
    if [ -f "$CATALINA_PID" ]; then
        if [ -s "$CATALINA_PID" ]; then
            if [ -r "$CATALINA_PID" ]; then
                PID=$(cat "$CATALINA_PID")
                ps -p $PID >/dev/null 2>&1
                if [ $? -eq 0 ] ; then
                    echo "Tomcat is running with PID $PID."
                    exit 0
                else
                    echo "Tomcat PID file exists with PID $PID but no process found."
                    exit 1
                fi
            else
                echo "Unable to read Tomcat PID file."
                exit 4
            fi
        else
            echo "Tomcat PID file is empty."
            exit 4
        fi
    else
        echo "No PID file found. Tomcat does not seem to be running."
        exit 3
    fi
    exit 4
}

case $1 in
start)
    tomcat_start
    ;;
stop)
    tomcat_stop
    ;;
restart)
    tomcat_stop
    tomcat_start
    ;;
status)
    tomcat_status
    ;;
esac
